# root/etc/s6-overlay/s6-rc.d/ambed/run
#!/command/with-contenv bash
#/sbin/rmmod ftdi_sio
#/sbin/rmmod usbserial
#/ambed/ambed 127.0.0.1
exec /ambed/ambed 127.0.0.1
