# root/etc/s6-overlay/s6-rc.d/ambed/run
#!/command/with-contenv sh

exec /ambed/ambed 127.0.0.1
